# network/cgc-utils/hypervisor-network.conf
#
description     "start hypervisor networking"

start on (floating-ip-bind or started nova-network)

task
exec /home/team/cgc-utils/configure_hypervisor.sh start-floating
