# network/cgc-utils/hypervisor-basic-network.conf
#
description     "start hypervisor basic networking"

start on started nova-network

task
exec /sbin/iptables -t nat -A POSTROUTING -o br100 -j MASQUERADE
