# network/cgc-utils/one_dot_thirty.conf
#
description     "start one dot thirty VM"

start on started qemu-kvm

task
exec /home/team/cgc-utils/boot-one_dot_thirty.sh start
