# network/gocap/cap.conf
#
description "Start cap."
start on started mysql
task
exec /usr/local/share/cap/cap.sh start
