# network/gocap/configure-mysql.conf
#
description "configure mysql in /run/"
start on starting mysql
task
exec /usr/local/share/mysql_default_schema/install_default_schema.sh
