# network/gocap/filter.conf
#
description "Start filter."
start on net-device-up IFACE=eth0
task
exec /usr/local/share/cap/filter.sh start
